# debian/c100-bootupd.service
[Unit]
Description=Update OKdo Nano C100 bootloader

[Service]
ExecStart=/usr/sbin/c100-bootupd

[Install]
WantedBy=multi-user.target
